# docker/s6-rc.d/afstart/up
#!/command/execlineb -P
/etc/ShellCrash/start.sh afstart
