# docker/s6-rc.d/crond/run
#!/command/execlineb -P
fdmove -c 2 1
exec crond -f -l 8