# docker/s6-rc.d/shellcrash/run
#!/bin/sh
set -e

. /etc/ShellCrash/configs/command.env

exec /bin/sh -c "${COMMAND} > /dev/null"
