# docker/s6-rc.d/bfstart/up
#!/command/execlineb -P
/etc/ShellCrash/start.sh bfstart
